# systemd/docker-cleanup.timer
[Unit]
Description=Cleanup timer for Docker images

[Timer]
OnUnitInactiveSec=12h

[Install]
WantedBy=timers.target
